# deploy/satoshistacker-alts.service
# Multi-asset alt-stacker (SOL/ETH/HYPE) on KRAKEN (USD pot). Sequential per-asset LLM
# calls, 4h decision loop, weekly per-asset technical self-tune, daily/weekly Telegram.
# Kraken has NO spot testnet, so this runs dry_run (paper on real Kraken prices). To go
# live (after the backtest gate + review): change to `--mode live` and set
# LIVE_TRADING_CONFIRMED=yes + KRAKEN_API_KEY/SECRET (trade-only, withdrawals disabled).
# Install: sudo cp deploy/satoshistacker-alts.service /etc/systemd/system/
#          sudo systemctl daemon-reload && sudo systemctl enable --now satoshistacker-alts
[Unit]
Description=SatoshiStacker multi-asset alt-stacker (SOL/ETH/HYPE)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=satoshi
WorkingDirectory=/opt/satoshistacker
EnvironmentFile=/opt/satoshistacker/.env
ExecStart=/opt/satoshistacker/.venv/bin/python -m agent.alts_main --mode dry_run
Restart=always
RestartSec=15
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal
SyslogIdentifier=satoshistacker-alts

[Install]
WantedBy=multi-user.target
